# etc/init.d/spice-webdavd-openrc
#!/sbin/openrc-run
# Copyright (c) 2007-2021 The Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts the SPICE WebDAV daemon for guest file sharing"

# Configuration is sourced automatically from the companion conf.d file
# /etc/conf.d/spice-webdavd

command="/usr/sbin/spice-webdavd"
command_args="-p 9843 ${spice_webdavd_opts}"
pidfile="/run/${RC_SVCNAME}.pid"

# This supervisor ensures the daemon is always running and restarts it if it crashes.
supervisor=supervise-daemon

depend() {
	# Needs D-Bus and the Avahi daemon to be running.
	need dbus avahi-daemon
}
