# debian/postinst
#!/usr/bin/env bash

# -- Add service files to runlevel/default.

rc-update add avahi-daemon-openrc default
rc-update add spice-webdavd-openrc default

# -- Update services list.

rc-update -u
